# etc/supervisor/supervisord.staging.conf
[supervisord]
nodaemon=true

[program:queue]
directory=/app
command=php artisan queue:work
autostart=yes
autorestart=true
